# backend/scripts/backup/crontab
# BoxMeOut Stella — Backup cron schedule
# Daily backup at 2:00 AM UTC
0 2 * * * /scripts/backup.sh >> /var/log/backup.log 2>&1
# Weekly verification on Sundays at 3:00 AM UTC
0 3 * * 0 /scripts/verify.sh >> /var/log/backup.log 2>&1
